# onchain_beast.service
[Unit]
Description=OnChain Beast - Side-Wallet Tracer (Solana)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=onchain
WorkingDirectory=/opt/onchain_beast
ExecStart=/opt/onchain_beast/target/release/onchain_beast
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
EnvironmentFile=/opt/onchain_beast/.env

[Install]
WantedBy=multi-user.target
